# dosbox.conf
[program:dosbox]
command=dosbox
autorestart=true
fullscreen=true

