# dosbox.conf
[program:dosbox]
command=dosbox
autorestart=true
fullscreen=true

[autoexec]
mount c /root/dos/civ
c:
civ.exe

