# dosbox.conf
[program:dosbox]
command=dosbox
autorestart=true
fullscreen=true
fulldouble=true
fullresolution=1366x768
windowresolution=1366x768


